# docker/tests/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=debug

[program:vpp]
command=/usr/bin/vpp -c /etc/vpp/startup.conf
autorestart=true
redirect_stderr=true
priority=1

[program:watchfrr]
command=/usr/lib/frr/watchfrr -d -F traditional zebra staticd
autorestart=true
redirect_stderr=true
priority=2

[program:netopeer2-server]
command=/usr/bin/netopeer2-server -d
autorestart=true
redirect_stderr=true
priority=3
